# daemon/codex-usage-daemon.service
[Unit]
Description=CodexMeter Usage Daemon (BLE)
Requires=bluetooth.target
After=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 DAEMON_PATH --transport serial
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
